# instahype.service
[Unit]
Description=Run instahype

[Service]
Type=simple
ExecStart=/usr/bin/instahype --cookies-db /var/cache/instahype/cookies.db --start-delay
