# instahype.timer
[Unit]
Description=Run instahype

[Timer]
OnBootSec=0
OnUnitActive=7600
AccuracySec=1

[Install]
WantedBy=timers.target
